# scripts/99-hdmi_sound.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/tobias/dotfiles/scripts/hdmi_suspend_toggle.sh"
